# Repository/rc.local
#!/bin/sh
set -e
#/etc/scripts/startup.sh
