# Repository/rc.local
#!/bin/sh
set -e
#/etc/scripts/startup.sh
sleep 15
#/etc/scripts/scheduled/avscan.sh
#/etc/scripts/scheduled/sync.sh
