# Repository/rc.local
#!/bin/sh
set -e
#/etc/scripts/startup.sh
sleep 15
#cpulimit -b -e clamscan -l 75 > /dev/null 2>&1
#/etc/scripts/scheduled/avscan.sh
#/etc/scripts/scheduled/sync.sh
