# Repository/rc.local
#!/bin/sh
#(sleep 120
#/etc/scripts/startup.sh)
#
