# Repository/rc.local
#!/bin/sh
set -e
#/etc/scripts/startup.sh &
#/etc/scripts/scheduled/avscan.sh &
#/etc/scripts/scheduled/sync.sh &
